# crontab
SHELL=/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin
59 23 * * * root . /etc/environment; /usr/local/bin/upload_to_s3.sh >/proc/1/fd/1 2>/proc/1/fd/2
